# scripts/bot-kt.service
[Unit]
Description=Reboots bot-kt if not running
[Service]
ExecStart=cd /home/mika/bot-kt && exec /home/mika/bot-kt/start.sh
