# scripts/bot-kt.service
[Unit]
Description=%hReboots bot-kt if not running
[Service]
ExecStart=/home/mika/bot-kt/start.sh
